# personal-brand-engine/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:api]
command=python -m uvicorn api.main:app --host 0.0.0.0 --port 8080
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api_error.log
environment=PYTHONPATH="/app"

[program:scheduler]
command=python -m scheduler.runner
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/scheduler.log
stderr_logfile=/app/logs/scheduler_error.log
environment=PYTHONPATH="/app"
